# contrib/retrieve-p2p-links-trigger.service
[Unit]
Description=Trigger retrieve-p2p-links regularly

[Service]
Type=oneshot
ExecStart=/bin/sh -c "cat /var/run/retrieve-p2p-links.pid | xargs -r kill -USR1"

[Install]
WantedBy=multi-user.target
